# systemd/coretex-corpus.service
[Unit]
Description=CoreTex retrieval corpus generation
Documentation=file:///root/cortex/docs/CORETEX_LAUNCH_PLAN_v2.md
After=network.target local-fs.target
Wants=network.target
StartLimitIntervalSec=600
StartLimitBurst=8

[Service]
Type=simple
User=root
WorkingDirectory=/root/cortex

# Pinned model + runtime config — matches calibration profile.
# Subprocess-level BLAS thread pinning lives in bi_encoder_runner.py
# (commit 4037885) so we don't redundantly set OMP/MKL/OPENBLAS here.
Environment=HF_HUB_CACHE=/var/lib/coretex/model-cache
Environment=HF_HUB_OFFLINE=1
Environment=CORETEX_BIENCODER=pinned
Environment=CORETEX_BIENCODER_MODE=streaming
Environment=CORETEX_BIENCODER_PYTHON=/root/cortex/.venv/bin/python
Environment=BIENCODER_NUM_THREADS=16
Environment=BIENCODER_INNER_BATCH=64

# Operator overrides go here (domains, seed count, output path).
# File is optional — defaults baked into run-corpus-generator.sh.
EnvironmentFile=-/etc/default/coretex-corpus

ExecStart=/root/cortex/scripts/run-corpus-generator.sh

# Auto-recovery: if the generator dies (NOT a clean exit), restart.
# A clean exit (success or explicit failure-no-retry) leaves the unit
# stopped so we don't loop on a corrupt corpus.
Restart=on-failure
RestartSec=15s

# Resource hygiene
LimitNOFILE=65536
# Make the corpus generator a LATE OOM-kill target. The kernel should
# kill anything else first.
OOMScoreAdjust=-500

# stdout/stderr → durable log file (in addition to journal).
# Append so successive restarts keep prior history.
StandardOutput=append:/var/lib/coretex/corpus-epoch-0-launch.log
StandardError=append:/var/lib/coretex/corpus-epoch-0-launch.log

# Kill the entire process tree on stop (bash→node→python).
# control-group ensures the python subprocess can't be left as zombie.
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=120

# Crucial: this unit runs in the SYSTEM scope, NOT user.slice.
# It is therefore IMMUNE to user-session lifecycle events
# (systemd-logind session-end, SSH disconnect, Claude harness restart).
# That is the entire reason this unit exists — see commit deec648
# context and CORETEX_PRODUCTION_RUNBOOK.md §"Crash recovery".
Slice=coretex.slice

[Install]
WantedBy=multi-user.target
